# Build/deb-meta/etc/systemd/system/mqttnet-server.service
[Unit]
Description=MQTTnetServer
After=network.target

[Service]
ExecStart=/opt/MQTTnetServer/MQTTnet.Server
WorkingDirectory=/opt/MQTTnetServer
StandardOutput=inherit
StandardError=inherit
Restart=always
User=mqttnetsrv

[Install]
WantedBy=multi-user.target